# .config/systemd/user/mount-drive.service
[Unit]
Description=Monta drives remotos (rclone) após a rede estar disponível
Documentation=https://rclone.org/commands/rclone_mount/
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/bin/rclone

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/mount_drive.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=default.target
